# packages/peppy-meter/peppymeter.service
[Unit]
Description=Peppy Meter
After=network-online.target

[Service]
Type=idle
Group=www-data
WorkingDirectory=/opt/peppymeter
ExecStart=/usr/bin/python3 /opt/peppymeter/peppymeter.py

[Install]
WantedBy=multi-user.target

