# packages/camillagui/camillagui.service
[Unit]
Description=CamillaDSP GUI
After=network-online.target

[Service]
Type=idle
Group=www-data
ExecStart=/usr/bin/python3 /opt/camillagui/main.py --config /opt/camillagui/config/camillagui.yml

[Install]
WantedBy=multi-user.target

