# packages/runonce/run_once.service
[Unit]
Description=The run_once service run scripts once during boot.
After=multi-user.target

[Service]
Type=simple
ExecStart=/bin/bash /usr/local/bin/runonce

[Install]
WantedBy=multi-user.target

